# whatsapp-bot.service
[Unit]
Description=WhatsApp UserBot v4.0 - AI Powered
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/whatsapp-bot
Environment="GEMINI_API_KEY=your-gemini-api-key-here"
Environment="HEADLESS=true"
ExecStart=/home/YOUR_USERNAME/whatsapp-bot/venv/bin/python /home/YOUR_USERNAME/whatsapp-bot/whatsapp_userbot.py
Restart=always
RestartSec=10
StandardOutput=append:/home/YOUR_USERNAME/whatsapp-bot/userbot.log
StandardError=append:/home/YOUR_USERNAME/whatsapp-bot/userbot.log

[Install]
WantedBy=multi-user.target
